# art-exhibits/com.jordan.art-exhibits-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jordan.art-exhibits-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jordan/jordan-os/art-exhibits/.venv/bin/python</string>
        <string>/Users/jordan/jordan-os/art-exhibits/main.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- Run on the 1st of every month at 9:00 AM -->
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/jordan/jordan-os/art-exhibits</string>

    <key>StandardOutPath</key>
    <string>/Users/jordan/jordan-os/art-exhibits/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jordan/jordan-os/art-exhibits/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/jordan/.local/bin</string>
    </dict>
</dict>
</plist>
